# obs-auto-moc/systemd/obs-auto-moc-pipeline.timer
[Unit]
Description=Periodic obs-auto-moc pipeline tick

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
